# internal/ignition/resources/target/systemd/dpf-ovs.service
[Unit]
Description=DPF OVS setup
After=network.target
ConditionPathExists=/etc/mlnx-release

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/dpf-ovs-script.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
