# internal/ignition/resources/target/systemd/dpu-agent.service
[Unit]
Description=DPF DPU Agent - Provisioning and Configuration
After=tmfifo-agent-link.service install-dpu-agent.service dpu-fw-upgrade.service
Before=nodeip-configuration.service kubelet-dependencies.target ovs-configuration.service
Requires=tmfifo-agent-link.service
Wants=install-dpu-agent.service dpu-fw-upgrade.service
ConditionPathExists=/etc/mlnx-release
ConditionPathExists=/usr/local/bin/dpu-agent

[Service]
Type=simple
EnvironmentFile=/etc/dpf/environment
ExecStart=/usr/local/bin/dpu-agent \
  --dpu-name $DPUName \
  --dpu-namespace $DPUNamespace \
  --dpu-uid $DPUUID \
  --dpuflavor /etc/dpf/dpuflavor.yaml \
  --skip-containerd-config \
  --skip-dns-config \
  --skip-kernel-cmd-line \
  --skip-network-config \
  --skip-remove-builtin-kubelet \
  --skip-configure-kubelet \
  --skip-start-kubelet \
  --skip-ovs-raw-script \
  --kubeadm-secret-name=unused
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
