# internal/ignition/resources/target/systemd/install-dpu-agent.service
[Unit]
Description=Install DPU agent
After=tmfifo-agent-link.service
Requires=tmfifo-agent-link.service
Before=machine-config-daemon-firstboot.service

[Service]
Type=oneshot
EnvironmentFile=/etc/dpf/environment
RemainAfterExit=yes
ExecStart=/usr/local/bin/install-dpu-agent.sh

[Install]
WantedBy=multi-user.target
