# internal/ignition/resources/target/systemd/pf-monitor.service
[Unit]
Description=DPF PF Link State Monitor
After=tmfifo-agent-link.service
Requires=tmfifo-agent-link.service

[Service]
Type=simple
EnvironmentFile=/etc/dpf/environment
ExecStart=/usr/local/bin/pf-monitor.sh
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
