# internal/ignition/resources/target/systemd/setup-vfs-devlink.service
[Unit]
Description=Setup Host VFs and Activate devlink
Before=install-dpu-agent.service
After=network-pre.target
# Don't start if this file exists:
ConditionPathExists=!/etc/mlnx-release

[Service]
EnvironmentFile=/etc/dpf/identity
Type=oneshot
ExecStart=/usr/local/bin/setup-vfs-devlink.sh

[Install]
WantedBy=multi-user.target
