# internal/ignition/resources/target/systemd/setup-vfs-devlink.service
[Unit]
Description=Setup Host VFs and Activate devlink
After=tmfifo-agent-link.service
Requires=tmfifo-agent-link.service
Before=install-dpu-agent.service
ConditionPathExists=!/etc/mlnx-release

[Service]
EnvironmentFile=/etc/dpf/environment
Type=oneshot
ExecStart=/usr/local/bin/setup-vfs-devlink.sh

[Install]
WantedBy=multi-user.target
